# shell/crontab.cron
MAIL=""
00 00 * * 0 /home/codestudy/shell/pushgit.sh
